# archiso/airootfs/etc/audit/rules.d/aegis.rules
# Aegis Linux - Audit Rules
# =============================================================================
# Comprehensive audit rules for security monitoring
# =============================================================================

# -----------------------------------------------------------------------------
# Clear existing rules and set buffer
# -----------------------------------------------------------------------------
-D
-b 8192
-f 1
--backlog_wait_time 60000

# -----------------------------------------------------------------------------
# Self-Auditing (audit the audit system)
# -----------------------------------------------------------------------------
-w /var/log/audit/ -p wa -k audit_log
-w /etc/audit/ -p wa -k audit_config
-w /etc/libaudit.conf -p wa -k audit_config
-w /etc/audisp/ -p wa -k audit_config

# -----------------------------------------------------------------------------
# Identity and Authentication Files
# -----------------------------------------------------------------------------
-w /etc/passwd -p wa -k identity
-w /etc/group -p wa -k identity
-w /etc/shadow -p wa -k identity
-w /etc/gshadow -p wa -k identity
-w /etc/security/opasswd -p wa -k identity

# PAM configuration
-w /etc/pam.d/ -p wa -k pam
-w /etc/security/ -p wa -k pam

# SSH configuration
-w /etc/ssh/sshd_config -p wa -k sshd
-w /etc/ssh/ssh_config -p wa -k ssh

# -----------------------------------------------------------------------------
# Privilege Escalation
# -----------------------------------------------------------------------------
-w /etc/sudoers -p wa -k sudoers
-w /etc/sudoers.d/ -p wa -k sudoers
-w /usr/bin/sudo -p x -k sudo_exec
-w /usr/bin/su -p x -k su_exec

# Monitor setuid/setgid program execution
-a always,exit -F arch=b64 -S execve -C uid!=euid -F euid=0 -k setuid
-a always,exit -F arch=b32 -S execve -C uid!=euid -F euid=0 -k setuid
-a always,exit -F arch=b64 -S execve -C gid!=egid -F egid=0 -k setgid
-a always,exit -F arch=b32 -S execve -C gid!=egid -F egid=0 -k setgid

# -----------------------------------------------------------------------------
# Login and Session Monitoring
# -----------------------------------------------------------------------------
-w /var/log/lastlog -p wa -k login
-w /var/log/faillog -p wa -k login
-w /var/log/btmp -p wa -k login
-w /var/log/wtmp -p wa -k login
-w /var/run/utmp -p wa -k login

# -----------------------------------------------------------------------------
# Cron and Scheduled Tasks
# -----------------------------------------------------------------------------
-w /etc/crontab -p wa -k cron
-w /etc/cron.d/ -p wa -k cron
-w /etc/cron.daily/ -p wa -k cron
-w /etc/cron.hourly/ -p wa -k cron
-w /etc/cron.monthly/ -p wa -k cron
-w /etc/cron.weekly/ -p wa -k cron
-w /var/spool/cron/ -p wa -k cron

# Systemd timers
-w /etc/systemd/system/*.timer -p wa -k systemd_timer

# -----------------------------------------------------------------------------
# Network Configuration
# -----------------------------------------------------------------------------
-w /etc/hosts -p wa -k network
-w /etc/hostname -p wa -k network
-w /etc/resolv.conf -p wa -k network
-w /etc/NetworkManager/ -p wa -k network

# Firewall configuration
-w /etc/nftables.conf -p wa -k firewall
-w /etc/iptables/ -p wa -k firewall

# -----------------------------------------------------------------------------
# Kernel and Module Operations
# -----------------------------------------------------------------------------
-w /sbin/insmod -p x -k modules
-w /sbin/rmmod -p x -k modules
-w /sbin/modprobe -p x -k modules
-w /etc/modprobe.d/ -p wa -k modules

# Monitor module loading/unloading
-a always,exit -F arch=b64 -S init_module,finit_module -k module_load
-a always,exit -F arch=b64 -S delete_module -k module_unload

# Kernel parameters
-w /etc/sysctl.conf -p wa -k sysctl
-w /etc/sysctl.d/ -p wa -k sysctl

# -----------------------------------------------------------------------------
# AppArmor
# -----------------------------------------------------------------------------
-w /etc/apparmor/ -p wa -k apparmor
-w /etc/apparmor.d/ -p wa -k apparmor
-w /sbin/apparmor_parser -p x -k apparmor

# -----------------------------------------------------------------------------
# Boot and Startup
# -----------------------------------------------------------------------------
-w /etc/default/grub -p wa -k boot
-w /etc/grub.d/ -p wa -k boot
-w /boot/ -p wa -k boot

# Systemd configuration
-w /etc/systemd/ -p wa -k systemd
-w /lib/systemd/ -p wa -k systemd

# -----------------------------------------------------------------------------
# Package Management
# -----------------------------------------------------------------------------
-w /usr/bin/pacman -p x -k package_mgmt
-w /etc/pacman.conf -p wa -k package_mgmt
-w /etc/pacman.d/ -p wa -k package_mgmt

# -----------------------------------------------------------------------------
# Time Changes
# -----------------------------------------------------------------------------
-a always,exit -F arch=b64 -S adjtimex,settimeofday,clock_settime -k time_change
-a always,exit -F arch=b32 -S adjtimex,settimeofday,clock_settime -k time_change
-w /etc/localtime -p wa -k time_change

# -----------------------------------------------------------------------------
# File System Mounts
# -----------------------------------------------------------------------------
-a always,exit -F arch=b64 -S mount,umount2 -k mount
-a always,exit -F arch=b32 -S mount,umount2 -k mount
-w /etc/fstab -p wa -k mount

# -----------------------------------------------------------------------------
# User/Group Changes
# -----------------------------------------------------------------------------
-w /usr/sbin/useradd -p x -k user_modification
-w /usr/sbin/userdel -p x -k user_modification
-w /usr/sbin/usermod -p x -k user_modification
-w /usr/sbin/groupadd -p x -k group_modification
-w /usr/sbin/groupdel -p x -k group_modification
-w /usr/sbin/groupmod -p x -k group_modification

# -----------------------------------------------------------------------------
# Suspicious Activity
# -----------------------------------------------------------------------------
# Reconnaissance tools
-w /usr/bin/nmap -p x -k recon
-w /usr/bin/tcpdump -p x -k recon
-w /usr/sbin/tcpdump -p x -k recon

# Power events
-w /usr/bin/systemctl -p x -k systemctl_exec

# Injection tools
-w /usr/bin/strace -p x -k injection
-w /usr/bin/ltrace -p x -k injection

# -----------------------------------------------------------------------------
# File Deletion
# -----------------------------------------------------------------------------
-a always,exit -F arch=b64 -S unlink,unlinkat,rename,renameat -F auid>=1000 -F auid!=-1 -k file_deletion
-a always,exit -F arch=b32 -S unlink,unlinkat,rename,renameat -F auid>=1000 -F auid!=-1 -k file_deletion

# -----------------------------------------------------------------------------
# Unauthorized Access Attempts
# -----------------------------------------------------------------------------
-a always,exit -F arch=b64 -S open,creat,truncate,ftruncate,openat -F exit=-EACCES -k access_denied
-a always,exit -F arch=b64 -S open,creat,truncate,ftruncate,openat -F exit=-EPERM -k access_denied

# -----------------------------------------------------------------------------
# Process ID Change
# -----------------------------------------------------------------------------
-a always,exit -F arch=b64 -S setuid,setgid,setreuid,setregid -k identity_change

# -----------------------------------------------------------------------------
# Make Configuration Immutable
# -----------------------------------------------------------------------------
# This MUST be the last rule - no changes after boot
-e 2
